# launchagent/de.nicx.phonebook-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Login-Autostart für Phonebook Server.

  Installation (Pfad ggf. auf /Applications anpassen):
      cp launchagent/de.nicx.phonebook-server.plist ~/Library/LaunchAgents/
      launchctl load -w ~/Library/LaunchAgents/de.nicx.phonebook-server.plist

  `open -a` ist Absicht: läuft die App schon, wird die bestehende Instanz aktiviert,
  statt eine zweite zu starten, die dann am belegten Port 8081 scheitern würde.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.nicx.phonebook-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/open</string>
        <string>-a</string>
        <string>/Applications/PhonebookServer.app</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
